# cron/crontab
*/1      *       *       *       *       root /home/amit/cronScript/updateHNFeed.sh
